# ntripCaster/checkbase2/check_caster.timer
[Unit]
Description=Check centipede base RTK and logs

[Timer]
OnCalendar=*-*-* *:*:0/15
AccuracySec=1s

[Install]
WantedBy=timers.target
